# scripts/com.savage.shc-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.savage.shc-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/robsavage/projects/savage-health-center/scripts/sync-shc.sh</string>
  </array>

  <!-- 4×/day: 7am, 1pm, 6pm, 10pm. Times the laptop is most likely awake.
       launchd coalesces missed runs (e.g. asleep at 7am) into the next wake. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>13</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>22</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/robsavage/Library/Logs/shc-sync/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/robsavage/Library/Logs/shc-sync/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
